# services/rotate-video/run
#!/bin/sh
exec setuidgid pi node /home/pi/teslacam/src/rotate-video.js
